# squashfs-root/etc/uci-defaults/luci-hwacc
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@hwacc[-1]
	add ucitrack hwacc
	set ucitrack.@hwacc[-1].init=hwacc
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
