# squashfs-root/etc/uci-defaults/luci-theme-xeye
#!/bin/sh
uci batch <<-EOF
	set luci.themes.LuciXEyE=/luci-static/LuciXEyE
	set luci.main.mediaurlbase="/luci-static/LuciXEyE"
        commit luci
EOF
	
