# squashfs-root/etc/hotplug.d/net/002-hwacc
#!/bin/sh
if [ -f /etc/init.d/hwacc ]  && [ "$DEVICENAME" = "br-lan" ]; then
      if [ "$ACTION" = "remove" ]; then
	/etc/init.d/hwacc stop
      else
	/etc/init.d/hwacc start
      fi;
fi


