# squashfs-root/usr/lib/opkg/info/bandwidth-pandorabox.list
/etc/config/bandwidth
/etc/init.d/bandwidth
/usr/bin/bw
