# squashfs-root/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/etc/banner.failsafe
/etc/resolv.conf
/lib/wifi/mt7612.sh
/etc_ro/wlan/MT7620_AP_2T2R-4L_V15.BIN
/etc/init.d/done
/bin/config_generate
/var
/sbin/config-vlan
/etc/hotplug.d/net/002-nas-acc
/lib/wifi/librt2860v2.sh
/etc/init.d/system
/etc/preinit
/lib/upgrade/platform.sh
/sbin/ndp_cfg
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/etc/hotplug.d/wifi/001-wifi.sh
/lib/preinit/10_indicate_preinit
/lib/functions/uci-defaults-new.sh
/sbin/esw_port_reset
/etc/config/system
/etc_ro/wlan/RT3662_AP_2T2R_V0_0.bin
/bin/cpu-affinity
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/shadow
/etc/rc.common
/etc/Wireless/SingleSKU_2G.dat
/etc/uci-defaults/10_migrate-shadow
/lib/netifd/wireless/mt7612.sh
/lib/preinit/10_indicate_failsafe
/etc_ro/wlan/RT3350_AP_1T1R_V1_0.bin
/init
/lib/netifd/wireless/rt2860v2.sh
/lib/preinit/80_mount_root
/etc_ro/wlan/RT3051_AP_1T2R_V1_0.bin
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc_ro/wlan/RT3052_AP_2T2R_V1_1.bin
/etc/init.d/led
/bin/wifitoggle
/etc/hotplug.d/block/20-usb_timecloud_led
/lib/functions/network.sh
/etc/init.d/sysctl
/bin/powersave
/sbin/hotplug-call
/etc/fw_env.config
/bin/wps_pbc
/etc/rc.button/failsafe
/etc_ro/wlan/MT7628_AP_2T2R-4L_V13.BIN
/etc/services
/lib/functions/system.sh
/bin/board_detect
/etc/hotplug.d/switch/001-switch.sh
/lib/wifi/rt2860v2.sh
/etc_ro/wlan/RT3050_AP_1T1R_V1_0.bin
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/hotplug.d/module/01-i2c
/etc_ro/wlan/RT5350_AP_1T1R_V1_0.bin
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/lib/preinit/40_run_failsafe_hook
/etc_ro/wlan/RT3352_AP_2T2R_V11.BIN
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/etc/functions.sh
/bin/smp-optimize
/sbin/firstboot
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/etc/uci-defaults/leds
/rom/note
/etc/inittab
/etc/openwrt_release
/etc_ro/wlan/MT7620_AP_2T2R-4L_V13.BIN
/sbin/wifi
/etc/profile
/etc/Wireless/SingleSKU_5G.dat
/bin/factory_system_load
/etc/uci-defaults/02_network
/etc_ro/wlan/MT7620_AP_2T2R-4L_external LNA_internal PA_V15.BIN
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/etc_ro/wlan/MT7662E2_EEPROM_20130527.bin
/etc/init.d/umount
/lib/functions/service.sh
/lib/ralink.sh
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/lib/preinit/02_default_set_state
/etc/mtab
/etc/init.d/boot
/lib/functions/preinit.sh
/bin/ipcalc.sh
/etc_ro/wlan/RT3883_AP_3T3R_V0_1.bin
/etc/diag.sh
/lib/preinit/03_preinit_do_ralink.sh
