# squashfs-root/usr/lib/opkg/info/dnsmasq-full.list
/etc/dnsmasq.conf
/etc/hotplug.d/iface/25-dnsmasq
/etc/config/dhcp
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
/usr/share/dnsmasq/trust-anchors.conf
