# squashfs-root/usr/lib/opkg/info/etherwake.list
/usr/bin/etherwake
/etc/config/etherwake
/etc/init.d/etherwake
