# squashfs-root/usr/lib/opkg/info/hd-idle.list
/usr/bin/hd-idle
/etc/config/hd-idle
/etc/init.d/hd-idle
