# squashfs-root/usr/lib/opkg/info/hwnat.list
/etc/hotplug.d/net/002-hwacc
/etc/init.d/hwacc
/etc/uci-defaults/99-hwacc
/bin/hwnat
/bin/ac
/bin/acl
/bin/mtr
