# squashfs-root/usr/lib/opkg/info/luci-app-ngrokc.list
/usr/lib/lua/luci/view/ngrokc/overview_value.htm
/usr/lib/lua/luci/model/cbi/ngrokc/detail.lua
/etc/uci-defaults/luci-app-ngrokc
/usr/lib/lua/luci/view/ngrokc/overview_enabled.htm
/usr/lib/lua/luci/view/ngrokc/ngrokc_script.htm
/usr/lib/lua/luci/controller/ngrokc.lua
/usr/lib/lua/luci/i18n/ngrokc.zh-cn.lmo
/usr/lib/lua/luci/model/cbi/ngrokc/overview.lua
/etc/init.d/ngrokc
/etc/config/ngrokc
