# squashfs-root/usr/lib/opkg/info/luci-app-chinadns.list
/usr/lib/lua/luci/i18n/chinadns.zh-cn.lmo
/etc/uci-defaults/luci-chinadns
/etc/init.d/chinadns
/etc/config/chinadns
/usr/lib/lua/luci/controller/chinadns.lua
/usr/lib/lua/luci/model/cbi/chinadns.lua
