# squashfs-root/usr/lib/opkg/info/luci-app-xware3.list
/etc/uci-defaults/luci-xware3
/usr/lib/lua/luci/model/cbi/xware3.lua
/usr/lib/lua/luci/controller/xware3.lua
/etc/init.d/xware3
/etc/config/xware3
