# squashfs-root/usr/lib/opkg/info/odhcpd.list
/usr/sbin/odhcpd
/etc/init.d/odhcpd
/etc/uci-defaults/odhcpd.defaults
/usr/sbin/odhcpd-update
