# squashfs-root/usr/lib/opkg/info/p910nd.list
/usr/sbin/p910nd
/etc/init.d/p910nd
/etc/config/p910nd
