# squashfs-root/usr/lib/opkg/info/pdnsd.list
/usr/bin/pdnsd-ctl
/usr/sbin/pdnsd
/etc/init.d/pdnsd
