# squashfs-root/usr/lib/opkg/info/relayd.list
/usr/sbin/relayd
/etc/init.d/relayd
/etc/hotplug.d/iface/30-relay
