# squashfs-root/usr/lib/opkg/info/rpcd.list
/etc/init.d/rpcd
/etc/config/rpcd
/sbin/rpcd
/usr/share/rpcd/acl.d/unauthenticated.json
