# squashfs-root/usr/lib/opkg/info/samba36-server.list
/etc/config/samba
/usr/sbin/nmbd
/etc/samba/smb.conf.template
/etc/samba/valid.dat
/etc/init.d/samba
/etc/samba/upcase.dat
/usr/sbin/smbpasswd
/etc/hotplug.d/block/20-smb
/usr/sbin/smbd
/usr/sbin/samba_multicall
/etc/samba/lowcase.dat
