# squashfs-root/usr/lib/opkg/info/transmission-pandorabox-daemon.list
/etc/init.d/transmission
/usr/bin/transmission-daemon
/etc/config/transmission
