# squashfs-root/usr/lib/opkg/info/usb-modeswitch.list
/sbin/usbmode
/etc/usb-mode.json
/etc/hotplug.d/usb/20-usb_mode
/etc/init.d/usbmode
