# squashfs-root/etc/hotplug.d/wifi/001-wifi.sh
#!/bin/sh
logger -t wifi   "ACTION=$ACTION  MAC=$MAC  SPEED=$SPEED SIGNAL=$SIGNAL"

